# files/etc/rc.local
sleep 10
/etc/init.d/dnscrypt-proxy restart
exit 0
